# launchd/com.jira-timer.idle-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jira-timer.idle-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.local/bin/jt-idle-monitor</string>
    </array>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/jira-timer/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/jira-timer/launchd.err.log</string>
</dict>
</plist>
